# services/acp-cli-v2/dev.ack.acp-provider.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.ack.acp-provider</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tjb/code/ack-protocol/services/acp-cli-v2/run-provider.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tjb/code/ack-protocol/services/acp-cli-v2</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/tjb/code/ack-protocol/services/acp-cli-v2/logs/provider.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tjb/code/ack-protocol/services/acp-cli-v2/logs/provider.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/tjb</string>
    </dict>
</dict>
</plist>
